# artix-cs2-arena-iso/profile/overlay/etc/xdg/autostart/gamemode.desktop
[Desktop Entry]
Type=Application
Name=GameMode Session Setup
Comment=Apply gaming-optimized settings for this session
Exec=/usr/local/bin/gaming-mode
Terminal=false
Hidden=false
X-GNOME-Autostart-enabled=true
